# supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
logfile_maxbytes=0

[program:fastapi]
command=bash -c "FASTAPI_PORT=${PORT:-8000}; echo \"FastAPI starting on port $FASTAPI_PORT\"; uvicorn main:app --host 0.0.0.0 --port $FASTAPI_PORT"
directory=/app/api
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0

[program:shopify]
command=bash -c "cd /app/shopify-app && npx prisma db push --skip-generate --accept-data-loss 2>/dev/null; echo 'Shopify Remix starting on port 3000'; PORT=3000 DATABASE_URL='file:/app/shopify-app/prisma/dev.sqlite' npm run start"
directory=/app/shopify-app
environment=NODE_ENV=production
autostart=true
autorestart=true
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
